# etc/profile.d/netboot-venv.sh
# Netboot virtualenv
export PATH="/opt/netboot-venv/bin:$PATH"
